# debian/cron.d
*/5 * * * *	www-data	php -d suhosin.session.encrypt=false /usr/share/igestis/modules/Roundcube/roundcubemail/plugins/fetchmail_rc/bash_cron.php > /dev/null 2&>1
