# debian/postinst
#!/bin/sh
# postinst script for igestis-roundcube
#
# see: dh_installdeb(1)

. /usr/share/debconf/confmodule


#DEBHELPER#


case "$1" in
    configure)

            if [ -e /etc/igestis/debian-db.php ] ; then

                # On récupère les infos de la base de données.
                dbname=$(cat /etc/igestis/debian-db.php  | grep "\$dbname" | cut -d\' -f2)
                dbuser=$(cat /etc/igestis/debian-db.php  | grep "\$dbuser" | cut -d\' -f2)
                dbpass=$(cat /etc/igestis/debian-db.php  | grep "\$dbpass" | cut -d\' -f2)

                # On vérifie que la structure de mise à jour automatique est bien présente.
                mysql -u$dbuser -p$dbpass $dbname -e 'CREATE TABLE IF NOT EXISTS `MYSQL_UPDATES` (`sql_file` varchar(50) collate utf8_unicode_ci NOT NULL,`date` datetime NOT NULL) ENGINE=MyISAM DEFAULT CHARSET=utf8 COLLATE=utf8_unicode_ci;'

                # On récupère et lance la liste des mises à jour.
		
                ls -v1 /usr/share/dbconfig-common/data/igestis-roundcube/install/ | while read i ; do
                    file=$(mysql -u$dbuser -p$dbpass $dbname -e "SELECT sql_file FROM MYSQL_UPDATES WHERE sql_file='igestis-roundcube-$i';")
                    if [ -z "$file" ] ; then
                        mysql -u$dbuser -p$dbpass $dbname < /usr/share/dbconfig-common/data/igestis-roundcube/install/$i || true
                        mysql -u$dbuser -p$dbpass $dbname -e "INSERT INTO MYSQL_UPDATES (sql_file, \`date\`) VALUES ('igestis-roundcube-$i', NOW());"
                    fi
                done

                ls -v1 /usr/share/dbconfig-common/data/igestis-roundcube/upgrade/mysql/ | while read i ; do
                    file=$(mysql -u$dbuser -p$dbpass $dbname -e "SELECT sql_file FROM MYSQL_UPDATES WHERE sql_file='igestis-roundcube-$i';")
                    if [ -z "$file" ] ; then
                        mysql -u$dbuser -p$dbpass $dbname < /usr/share/dbconfig-common/data/igestis-roundcube/upgrade/mysql/$i || true
                        mysql -u$dbuser -p$dbpass $dbname -e "INSERT INTO MYSQL_UPDATES (sql_file, \`date\`) VALUES ('igestis-roundcube-$i', NOW());"
                    fi
                done

            fi


        webserver=apache2

        if [ ! -z "$(file /etc/apache2/conf.d/roundcube | grep "ASCII text, ")" ] ; then
            rm -f /etc/apache2/conf.d/igestis
        fi


        if [ -d /etc/$webserver/conf.d ] && [ ! -e /etc/$webserver/conf.d/roundcube ]; then
            ln -sf /etc/igestis/apache-roundcube.conf /etc/$webserver/conf.d/roundcube
        fi
        invoke-rc.d apache2 reload 3>/dev/null || true


        CONFFILE=/etc/igestis/config.php

    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0


