# debian/postinst
#!/bin/sh
# postinst script for igestis
#
# see: dh_installdeb(1)




webserver=apache2

#DEBHELPER#



case "$1" in
    configure)
    
    if [ ! -z "$(file /etc/apache2/conf.d/roundcube | grep "ASCII text")" ] ; then
		rm -f /etc/apache2/conf.d/roundcube
	fi
		
    
    if [ -d /etc/$webserver/conf.d ] && [ ! -e /etc/$webserver/conf.d/roundcube ]; then
               ln -sf /etc/igestis/apache-roundcube.conf /etc/$webserver/conf.d/roundcube
	fi
    invoke-rc.d apache2 reload 3>/dev/null || true

	dbname=$(cat /etc/igestis/debian-db.php  | grep "\$dbname" | cut -d\' -f2)
	dbuser=$(cat /etc/igestis/debian-db.php  | grep "\$dbuser" | cut -d\' -f2)
	dbpass=$(cat /etc/igestis/debian-db.php  | grep "\$dbpass" | cut -d\' -f2)

	mysql -u$dbuser -p$dbpass $dbname < /usr/share/dbconfig-common/data/igestis-roundcube/install/mysql 2> /dev/null || true

	find /usr/share/dbconfig-common/data/igestis-roundcube/upgrade/mysql/* 2> /dev/null | while read i ; do
		mysql -u$dbuser -p$dbpass $dbname < $i 2> /dev/null || true
	done
	

    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0


